# systemd/i3ds.target
[Unit]
Description=Target unit for all I3DS services

[Install]
WantedBy=multi-user.target
